# rootfs/etc/services.d/coturn/run
#!/usr/bin/with-contenv bash
turnserver -v
